# cookbooks/chatsecure_web/templates/default/upstart.conf.erb
# Generated by Chef, do not modify!
description "upstart configuration for chatsecure_web"

setuid <%= @service_user %>

start on runlevel [2345]
stop on runlevel [06]

respawn

script
    exec <%= @virtualenv_path %>/bin/uwsgi --chdir=<%= @app_root %>/<%= @app_name %> \
    --module=<%= @app_name %>.wsgi:application \
    --env DJANGO_SETTINGS_MODULE=<%= @app_name %>.settings \
    --master --pidfile=/tmp/<%= @service_name %>.pid \
    --die-on-term \
    --socket=127.0.0.1:<%= @app_port %> \
    --processes=<%= @app_workers %> \
    --uid=<%= @service_user_id %> --gid=<%= @service_user_gid %> \
    --harakiri=<%= @harakiri %> \
    --max-requests=<%= @max_requests %> \
    --vacuum \
    --home=<%= @virtualenv_path %> \
end script